# gcloud-workstation-vscode/systemd-assets/etc/systemd/system/docker.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd --data-root /home/user/docker -H fd:// --containerd=/run/containerd/containerd.sock
