# gcloud-workstation-phpstorm/systemd-assets/etc/systemd/system/openvpn-client@.service.d/override.conf
[Service]
WorkingDirectory=/home/user/vpn
ExecStart=
ExecStart=/usr/sbin/openvpn --suppress-timestamps --nobind --user nobody --config %i.ovpn
ProtectHome=false
